# sentry-telegram.service
[Unit]
Description=Sentry -> Telegram notifier
# Wait for the network so polling / Telegram API calls work on boot.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Project location. NOTE: systemd only expands ${PROJECT_PATH} in ExecStart;
# WorkingDirectory / EnvironmentFile / ReadWritePaths do NOT expand variables,
# so if you change this path you must update those literal paths below to match.
Environment=PROJECT_PATH=/opt/sentry-telegram

# Run as an unprivileged user. Create it once with:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin sentry-telegram
#   sudo chown -R sentry-telegram:sentry-telegram /opt/sentry-telegram
User=sentry-telegram
Group=sentry-telegram

# Project lives here; .env and state.db are read/written relative to this dir.
# (Keep in sync with PROJECT_PATH above.)
WorkingDirectory=/opt/sentry-telegram
# config.py reads .env from WorkingDirectory; this is a redundant safety net.
EnvironmentFile=-/opt/sentry-telegram/.env

# Use the project's venv interpreter; main.py resolves against WorkingDirectory.
ExecStart=${PROJECT_PATH}/venv/bin/python main.py

# Restart on crash, with a small backoff.
Restart=on-failure
RestartSec=5

# Send stdout/stderr to the journal (view with: journalctl -u sentry-telegram -f).
StandardOutput=journal
StandardError=journal

# --- hardening (relax if something legitimately needs more access) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# state.db lives in WorkingDirectory, so that path must stay writable.
ReadWritePaths=/opt/sentry-telegram

[Install]
WantedBy=multi-user.target
